# src/crontab
@reboot       /usr/bin/flock -n .sync.lock /sync.sh
*/15 * * * *  /usr/bin/flock -n .sync.lock /sync.sh
